# make.d/incus/systemd/rke2-server.service.d/10-manifests-install.conf
[Unit]
Wants=rke2-runtime-manifests-install.service
Wants=rke2-storage-manifests-install.service
Wants=rke2-networking-manifests-install.service
Wants=rke2-mesh-manifests-install.service
Wants=rke2-tekton-pipelines-manifests-install.service