# make.d/incus/systemd/rke2-server.service.d/10-manifests-install.conf
[Unit]
Wants=rke2-networking-manifests-install.service
Before=rke2-server.service