# make.d/incus/systemd/rke2-server.service.d/cilium-operator-scaling.conf
[Service]
ExecStartPost=/bin/sh -xc '/usr/local/sbin/rke2-cilium-operator-scaling'
