# make.d/incus/systemd/rke2-agent.service.d/kubeconfig.conf
[Service]
ExecStartPost=/bin/sh -xc '/srv/host/rke2/scripts/rke2-tools-configuration-directories.sh'
ExecStartPost=/bin/sh -xc '/srv/host/rke2/scripts/rke2-vip-kubeconfig.sh'
